# deploy/vpnctl-controller.service
[Unit]
Description=vpnctl controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vpnctl controller init --config /etc/vpnctl/controller.yaml
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
